# eol-first-boot/DEBIAN/postinst
#!/bin/bash
# postinst script for first-boot
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure|triggered)
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

if [ "$1" = configure ]; then
    # this script runs once, deleting itself and this package.
    update-rc.d eol-first-boot defaults

    # cleanup log files and such, in preparation for creating an
    # image of this system.

    service rsyslog stop
    service chrony stop

    cd /var/log

    for f in auth.log btmp daemon.log debug dmesg* kern.log lastlog messages syslog user.log wtmp; do
	rm -f ${f}.*
	cat /dev/null > $f
    done

    rm -f chrony/*

    cd -

    rm -f /etc/udev/rules.d/70-persistent-net.rules

    history -c
    rm -f /root/.bash_history /root/.viminfo

fi

exit 0
