# isfs-daq/DEBIAN/postinst
#!/bin/sh
# postinst script for isfs-daq
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

if [ "$1" = configure ]; then
    chown -R daq.eol /home/daq

    update-rc.d emerald defaults

    # edit /etc/default/gpsd:
    # DEVICES="/dev/ttyS3"
    # GPSD_OPTIONS="-b -n"

    # add /etc/chrony/chrony.conf
    # SHM1 from gpsd (if present) is from the kernel PPS_LDISC
    # module.  It includes PPS and will be accurate to a few ns
    # refclock SHM 1 offset 0.0 delay 0.0 refid PPS

    # add to /etc/modules
    # pps_ldisc


exit 0
